# packaging/debian/postinst
#!/usr/bin/env bash
set -euo pipefail

systemctl daemon-reload >/dev/null 2>&1 || true
systemctl enable linuxio.target >/dev/null 2>&1 || true
systemctl start linuxio.target >/dev/null 2>&1 || true

systemd-tmpfiles --create /usr/lib/tmpfiles.d/linuxio.conf >/dev/null 2>&1 || true

if [[ -d /etc/motd.d ]]; then
  ln -sf ../../run/linuxio/issue /etc/motd.d/linuxio >/dev/null 2>&1 || true
fi

exit 0
