# packaging/systemd/linuxio-bridge-socket-user.service
[Unit]
Description=Dynamic user for LinuxIO bridge sockets
Documentation=https://github.com/mordilloSan/LinuxIO
# PartOf provides stop propagation without circular dependency
# (when linuxio.service stops, this unit stops too)
PartOf=linuxio.service

[Service]
DynamicUser=yes
User=linuxio-bridge-socket
Group=linuxio-bridge-socket
Type=oneshot
ExecStart=/bin/true
RemainAfterExit=yes