# packaging/systemd/linuxio-bridge-socket-user.service
[Unit]
Description=Dynamic user for LinuxIO bridge sockets
Documentation=https://github.com/mordilloSan/LinuxIO
DefaultDependencies=no
Before=sockets.target
PartOf=linuxio-auth.socket

[Service]
DynamicUser=yes
User=linuxio-bridge-socket
Group=linuxio-bridge-socket
Type=oneshot
ExecStart=/bin/true
RemainAfterExit=yes
