# packaging/systemd/linuxio-auth.service
[Unit]
Description=LinuxIO Authentication Service
Documentation=https://github.com/mordilloSan/LinuxIO
Requires=linuxio-auth.socket
After=linuxio-auth.socket

[Service]
ExecStart=/usr/local/bin/linuxio-auth-helper
Restart=on-failure
RestartSec=5s

# Runs as root (required for PAM authentication)
# Minimal hardening like cockpit-session@.service because:
# - Child bridge processes inherit restrictions
# - Bridge needs full network/system access for WireGuard, etc.
# - Security comes from: PAM auth, per-user bridge isolation, socket permissions

[Install]
WantedBy=multi-user.target
